# scripts/upstart/codedeploy.conf
# upstart configuration

description "golang server"
author "mirco zeiss"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /home/ubuntu/gocode/bin/codedeploy
